# meta-m2/recipes-apps/custom-app/files/emb_znsd.service
[Unit]
Description=Embeded ZNSD Monitor

[Service]
Type=idle
StandardOutput=syslog
StandardError=syslog
ExecStart=/usr/bin/emb_znsd

[Install]
WantedBy=multi-user.target